# files/system/usr/lib/systemd/user/brew-upgrade.service
[Unit]
Description=Topgrade

[Service]
Type=oneshot
ExecStart=/home/linuxbrew/.linuxbrew/bin/topgrade --yes --cleanup
