# files/system/usr/lib/systemd/user/brew-upgrade.service
[Unit]
Description=Brew upgrade

[Service]
Type=oneshot
ExecStart=/home/linuxbrew/.linuxbrew/bin/brew upgrade
