# files/system/usr/lib/systemd/user/topgrade.timer
[Unit]
Description=Weekly topgrade

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
